# docker/picoclaw-container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:picoclaw-gateway]
command=su - picoclaw -s /bin/bash -lc "PICOCLAW_HOME=/home/picoclaw/.picoclaw PICOCLAW_CONFIG=/home/picoclaw/.picoclaw/picoclaw.json picoclaw gateway"
autorestart=true
priority=10
stdout_logfile=/var/log/picoclaw-gateway.log
stderr_logfile=/var/log/picoclaw-gateway.err
