# docker/openclaw-container/supervisord.conf
; Supervisor configuration for OpenClaw container
; Manages: openclaw-gateway, mcp-auth-proxy, mcporter

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; Increase startup timeout for slow-starting services
startsecs=5
startretries=3

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ═══════════════════════════════════════════════════════════════════════════════
; nginx — forwards 0.0.0.0:18789 → 127.0.0.1:18788 (openclaw-gateway loopback)
; ═══════════════════════════════════════════════════════════════════════════════
[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ═══════════════════════════════════════════════════════════════════════════════
; OpenClaw Gateway
; ═══════════════════════════════════════════════════════════════════════════════
[program:openclaw-gateway]
command=openclaw gateway run
user=openclaw
environment=HOME="/home/openclaw",USER="openclaw",NODE_OPTIONS="--max-old-space-size=768"
directory=/home/openclaw
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stopsignal=TERM

; Logging
stdout_logfile=/var/log/supervisor/openclaw-gateway.log
stderr_logfile=/var/log/supervisor/openclaw-gateway.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile_backups=3

; ═══════════════════════════════════════════════════════════════════════════════
; Bootstrap pairing approval
; Approves the first remote pairing request during the initial startup window.
; This is intended for the SuprClaw middleware, which is the only upstream device
; that should connect directly to the container gateway.
; ═══════════════════════════════════════════════════════════════════════════════
[program:pairing-bootstrap]
command=/usr/local/bin/auto-approve-pairing.sh
user=root
autostart=true
autorestart=false
startsecs=0
startretries=0
stopwaitsecs=5
stdout_logfile=/var/log/supervisor/pairing-bootstrap.log
stderr_logfile=/var/log/supervisor/pairing-bootstrap.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ═══════════════════════════════════════════════════════════════════════════════
; MCP Auth Proxy
; Routes MCP tool calls with authentication
; ═══════════════════════════════════════════════════════════════════════════════
[program:mcp-auth-proxy]
command=node /usr/local/bin/mcp-auth-proxy.js
user=root
directory=/home/openclaw
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM

; Logging
stdout_logfile=/var/log/supervisor/mcp-auth-proxy.log
stderr_logfile=/var/log/supervisor/mcp-auth-proxy.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile_backups=3

; ═══════════════════════════════════════════════════════════════════════════════
; Group: All OpenClaw services
; ═══════════════════════════════════════════════════════════════════════════════
[group:openclaw]
programs=nginx,openclaw-gateway,pairing-bootstrap,mcp-auth-proxy
priority=999
